# launchd/com.samuel.opencode-serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.samuel.opencode-serve</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/opencode</string>
    <string>serve</string>
    <string>--port</string>
    <string>4096</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/samuel.sanchez-moreno/Library/Logs/opencode-serve.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/samuel.sanchez-moreno/Library/Logs/opencode-serve.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/samuel.sanchez-moreno</string>
  </dict>
</dict>
</plist>
